# scripts/conf/gunicorn.conf
[program:gunicorn]
command=/srv/funny/pycon/scripts/gunicorn/start_non_detach.sh
user=funny
autostart=true
autorestart=true
redirect_stderr=true
